# crontab
# Host cron template (install on the EC2 host, not inside Docker).
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Every day at 15:30 UTC, finds the container by docker-compose
# service label ("spdx"), and runs the media cleanup job.
30 15 * * * c=$(docker ps -qf "label=com.docker.compose.service=spdx" | head -n 1) && [ -n "$c" ] && docker exec "$c" /usr/local/bin/python /spdx/src/manage.py cleanup_media
